# subscriber/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:subsciber]
command=php /usr/local/sbin/subscriber/subscriber-worker.php -h eventbus -p 6379 foo bar
stdout_logfile=/var/log/worker.log
autostart=true
autorestart=true
# user=www-data
stopsignal=KILL
numprocs=1